# ops/systemd/mc-paper.service.d/30-restart-sane.conf
[Service]
Restart=on-failure
RestartSec=5

ExecStop=
KillSignal=SIGTERM
TimeoutStopSec=60

SuccessExitStatus=143
