# ops/systemd/hybrid-engine.service.d/30-restart-exit12.conf
[Service]
Restart=on-failure
RestartPreventExitStatus=12
# Optional: treat 12 as clean for status dashboards
# SuccessExitStatus=12
